# systemd/Calc.service
[Unit]
Description=Dbus based calculator service
After=dbus.service

[Service]
Type=dbus
BusName=com.mateuszmidor.Calc
ExecStart=/usr/local/bin/CalcServer
Environment=COMMONAPI_DBUS_CONFIG=/usr/local/etc/commonapi-dbus/Calc-stub.ini
Restart=on-failure

[Install]
WantedBy=multi-user.target